# .config/systemd/user/container-woodpecker-agent.service
# container-woodpecker-agent.service
# autogenerated by Podman 4.7.0
# Sun Oct 29 13:12:20 PST 2023

[Unit]
Description=Podman container-woodpecker-agent.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	-d \
	--replace \
	--name=woodpecker-agent \
	--security-opt label=disable \
	--label=io.containers.autoupdate=registry \
	--volume=/media/backup/woodpecker:/etc/woodpecker:Z \
	--volume=/run/user/1000/podman/podman.sock:/var/run/podman.sock:Z \
	--env=DOCKER_HOST=unix:///var/run/podman.sock \
	--env=WOODPECKER_SERVER=woodpecker.samson.ph:9000 \
	--env=WOODPECKER_HOSTNAME=bomuhei \
	--env=WOODPECKER_AGENT_SECRET=gPxLm9rjB3ucYOXQwDHQN54Yq19ViJsZVYTnFXBhVxbtg2j63kHwP3BOIQ2F5zqO1HAfef1WxpPwQ7hIwH5h8w7wuLKlS7uzTFthrNPzKaTtloiTQl7vSDmbZmQlWdqZ docker.io/woodpeckerci/woodpecker-agent:latest
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
