# crontab
# Purge ~/Downloads
12 12 * * *     $HOME/.local/bin/purgeOld $HOME/Downloads

